# root/etc/init.d/serverchan
#!/bin/sh /etc/rc.common

USE_PROCD=1
START=99

start_service() {
        procd_open_instance
        procd_set_param command /usr/share/serverchan/serverchan
        procd_close_instance
        echo "ServerChan is starting now ..."
}

stop_service() {
        echo "Terminating serverchan process..."
}

service_triggers() {
    procd_add_reload_trigger serverchan
}

restart_service() {
	stop
        sleep 1
        start
}

