# root/etc/init.d/serverchan
#!/bin/sh /etc/rc.common
START=99
STOP=10


start() {
sleep 3
/usr/bin/serverchan/serverchan &
}
stop() {
kill -9 `pgrep -f serverchan`
}
