# root/etc/init.d/serverchan
#!/bin/sh /etc/rc.common
START=99
STOP=10


start() {
/usr/bin/serverchan/serverchan &
}
stop() {
kill -9 `pgrep -f serverchan` >/dev/null 2>&1
return 0
}
